# installers/debian/WFTableauCollector-0.1.0_amd64/DEBIAN/postinst
echo "*******************************************************"
echo " wftableaucollector is installed"
echo " YOU NEED TO CHANGE THE PROPERTY FILES LOCATED AT"
echo " /usr/local/wftableaucollector/config"
echo " "
echo " Then you can start the service using"
echo " systemctl enable wftableaucollector"
echo " systemctl start wftableaucollector"
echo " "
echo " Logs are located in"
echo " /usr/local/wftableacollector/log"
echo "*******************************************************"
